# dev/cloud/mkosi.extra/etc/systemd/system/malmo-metadata-firewall.service
# Hosted cloud image — block app-container egress to the cloud metadata endpoint
# (#251). Loads the standing nftables forward-hook DROP for 169.254.169.254
# (/etc/malmo/metadata-firewall.nft) so an untrusted app container cannot read the
# first-boot seed secrets via the classic cloud-metadata SSRF. ENVIRONMENT.md
# # Provisioning & first-boot.
#
# Ordered Before=docker.service and host-agent.service: the rule must be in place
# before docker starts any container and before host-agent launches the brain, so
# no app/control-plane container can ever reach 169.254.169.254 before the block.
# The rule sits in its own nftables table and only touches the forward hook, so it
# is independent of docker's own ruleset and of the host-root seed fetch (OUTPUT).
#
# Enabled by dev/cloud/mkosi.postinst.chroot (a .wants symlink) for both the lean
# production image and the boot-proof test lane that Include=..'s it.
[Unit]
Description=malmo — block container egress to the cloud metadata endpoint (169.254.169.254)
Documentation=https://github.com/malmoos/malmo/blob/main/docs/specs/ENVIRONMENT.md
Before=docker.service host-agent.service
After=nftables.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/sbin/nft -f /etc/malmo/metadata-firewall.nft

[Install]
WantedBy=multi-user.target
