# dist/systemd/dropins/smbd.service.d/molma.conf
# molma drop-in for smbd.service.
# Per BOOT.md # What downstream services do: smbd comes up even on partial
# storage failure. When the brain has a blocks_writes health issue active, it
# calls host-agent to reload Samba with a read-only share config — that
# reload path is separate from this drop-in. Read access stays available
# regardless.
[Unit]
After=molma-storage-ready.target
Wants=molma-storage-ready.target
