# dev/cloud/mkosi.extra/etc/systemd/system/docker.service.d/10-malmo-grow-root.conf
# Fail-closed on the boot-time root grow. malmo-grow-root.service extends the
# baked 8 GiB root partition + ext4 to fill the whole provider disk before any
# write-heavy service starts. docker image storage shares that one volume with the
# brain's SQLite store, so if the grow failed the root is still 8 GiB and a single
# app install can fill it and 500 login. Requires= refuses to start docker on an
# un-grown root rather than run into that. (The grow unit's own Before=docker.service
# already orders the grow ahead of us; this only adds the hard dependency.)
[Unit]
Requires=malmo-grow-root.service
