# dist/systemd/dropins/avahi-daemon.service.d/molma.conf
# molma drop-in for avahi-daemon.service.
# Per BOOT.md: Avahi is OFF the critical path for dashboard reachability.
# It depends on network.target (not network-online.target) and starts as soon
# as a link is up, so molma.local resolves before NTP, storage, Docker, or
# the brain. We add only the storage-target ordering for predictable boot
# sequencing — Wants=, not Requires=, so an Avahi crash does not affect
# storage assembly and a storage failure does not prevent name resolution.
[Unit]
After=molma-storage-ready.target
Wants=molma-storage-ready.target
