# dev/cloud/test/malmo-cloud-assertions.service
# Cloud boot-proof self-check (C2, #205; seed/gate scenarios C3a cloud-lane, #220).
# A boot oneshot that polls the control plane up and writes its verdict to the
# serial console, where dev/cloud/run-cloud-tests.sh greps it. Serial-driven (not
# SSH-driven) because hosted v1 ships no sshd (ENVIRONMENT.md # Access & files);
# the check runs in-guest over localhost, the real control-plane-up proof.
#
# Runs on EVERY boot (no run-once marker): #220 drives a multi-boot sequence over
# one persisted disk (un-seeded 503 → seeded gate → frozen-identity reboot), so
# the self-check must re-run each boot. It reads which scenario to assert from the
# malmo.assert credential (ImportCredential below; default un-seeded). This unit
# exists only in the boot-proof test image, never the lean production cut.
[Unit]
Description=malmo cloud boot-proof self-check (C2/C3a)
After=docker.service host-agent.service malmo-load-images.service
Wants=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
# The scenario directive (un-seeded | seeded | frozen:<box-id>) the harness sets
# per boot over SMBIOS; absent ⇒ the script defaults to the un-seeded 503 check.
ImportCredential=malmo.assert
# The script polls the control plane up (docker load + brain bootstrap + compose
# up race this unit) and can run for minutes; the default 90s start timeout would
# kill it mid-poll before it writes a verdict. No timeout — the QEMU harness owns
# the outer deadline.
TimeoutStartSec=0
ExecStart=/usr/local/bin/cloud-assertions.sh

[Install]
WantedBy=multi-user.target
