# dist/systemd/dropins/docker.service.d/molma.conf
# molma drop-in for docker.service.
# Per BOOT.md # ordering: Wants=, not Requires=, so a partial-storage failure
# does not block Docker (and therefore does not block host-agent / brain).
# The brain detects docker-data-root-misplaced via a separate health check
# when the bind isn't where it should be.
[Unit]
After=molma-storage-ready.target
Wants=molma-storage-ready.target
