# dist/systemd/molma-storage-ready.target
[Unit]
Description=molma storage assembly attempted
Documentation=https://github.com/molmaos/molma/blob/main/docs/specs/BOOT.md
# Best-effort milestone, not a strict gate. Per BOOT.md # The storage-ready
# target: services depending on storage Want= this target, not Require= it.
# Individual mount/verify failures do not stop downstream services starting.
After=local-fs.target
Wants=molma-storage-verify.service
