# dist/systemd/molma-storage-verify.service
[Unit]
Description=molma storage canary + enrollment-marker reporter
Documentation=https://github.com/molmaos/molma/blob/main/docs/specs/BOOT.md
# Reporter, not a gate. No OnFailure= per BOOT.md # Failure → recovery target.
# The reporter writes findings to /run/molma/health/storage.json which
# host-agent reads and forwards to the brain via GET /v1/health/storage.
#
# Note: the molma-storage-ready.target carries Wants=molma-storage-verify.service
# directly, so this unit needs no [Install] section — it's pulled in
# unconditionally whenever the target is reached. (The alternative,
# WantedBy=molma-storage-ready.target in [Install], would create a redundant
# enable-time symlink and a confusing "systemctl disable still leaves us
# wanted" state.)
DefaultDependencies=no
After=local-fs.target
Before=molma-storage-ready.target

[Service]
Type=oneshot
RemainAfterExit=yes
# Per BOOT.md # Hang protection on the storage chain: a hang trips this
# timeout; the chain moves on. A future enhancement (tracked in 0019 'what's
# next') is for a wrapper to record a storage-verify-timeout finding when
# systemd kills us here.
TimeoutStartSec=60s
ExecStart=/usr/lib/molma/molma-storage-verify
# /run/molma/health/ is created by the binary; tmpfiles.d would be cleaner
# (tracked in 0019 'what's next').
