# dev/cloud/mkosi.extra/etc/systemd/system/host-agent.service.d/10-malmo-grow-root.conf
# Fail-closed on the boot-time root grow (hosted image only). host-agent launches
# the brain, which writes the SQLite store on the shared root volume. If
# malmo-grow-root.service failed, the root is still the baked 8 GiB and a single
# app install can fill that shared docker+SQLite volume and 500 login — so refuse
# to start on an un-grown root rather than run into it. This lives in the tracked
# hosted ExtraTree (mkosi.extra/), NOT in the shared dist/systemd/host-agent.service
# — the appliance has no malmo-grow-root.service, so the Requires= must be scoped to
# the hosted profile via a drop-in. malmo-grow-root.service's own Before=host-agent.service
# supplies the ordering; this only adds the hard dependency.
[Unit]
Requires=malmo-grow-root.service
After=malmo-grow-root.service
